# .docker/config/sp-php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm8.2 --fpm-config %(GITPOD_REPO_ROOT)s/.docker/config/php-fpm.conf
process_name=%(program_name)s
priority=1001
autostart=true
startretries=3
autorestart=true
user=gitpod
